# dot-files/systemd-user/personal-audit.timer
[Unit]
Description=Run personal-audit.service

[Timer]
OnCalendar=hourly

[Install]
WantedBy=user-graphical-login.target
