# dot-files/systemd-user/personal-audit.service
[Unit]
Description=Audit my files for actions needed

[Service]
Type=oneshot
ExecStart=/home/bexelbie/bin/personal-audit
